# deploy/systemd/maintenance/musubi-retention.service
[Unit]
Description=Musubi Storage Retention Worker
After=network.target

[Service]
Type=oneshot
User=musubi
Group=musubi
WorkingDirectory=/opt/musubi/src
EnvironmentFile=/etc/musubi/.env
# Replace with actual command to run retention
ExecStart=/opt/musubi/.venv/bin/python -m musubi.ops.retention_cmd
