# deploy/systemd/musubi-lifecycle-worker.service
[Unit]
Description=Musubi lifecycle worker
Requires=docker.service
After=network-online.target docker.service musubi-api.service
Wants=network-online.target

[Service]
Type=simple
User=musubi
Group=musubi
WorkingDirectory=/opt/musubi/current
EnvironmentFile=/etc/musubi/.env.production
ExecStart=/usr/bin/env uv run python -m musubi.lifecycle.scheduler
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=musubi-lifecycle-worker

[Install]
WantedBy=multi-user.target
